# Enter/pkgs/conda-25.11.1-py313h06a4308_0/etc/profile.d/conda.csh
setenv CONDA_EXE "/home/task_176556791054345/conda-bld/conda_1765567979493/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_p/bin/conda";
setenv _CONDA_EXE "/home/task_176556791054345/conda-bld/conda_1765567979493/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_p/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/home/task_176556791054345/conda-bld/conda_1765567979493/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_p/bin/python";
setenv _CONDA_ROOT "/home/task_176556791054345/conda-bld/conda_1765567979493/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_p";

source "/home/task_176556791054345/conda-bld/conda_1765567979493/_h_env_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_placehold_p/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"