# Enter/etc/profile.d/conda.csh
setenv CONDA_EXE "/workspaces/all-in-rag/Enter/bin/conda";
setenv _CONDA_EXE "/workspaces/all-in-rag/Enter/bin/conda";
unsetenv _CE_M;
unsetenv _CE_CONDA;
setenv CONDA_PYTHON_EXE "/workspaces/all-in-rag/Enter/bin/python";
setenv _CONDA_ROOT "/workspaces/all-in-rag/Enter";

source "/workspaces/all-in-rag/Enter/lib/python3.13/site-packages/conda/shell/etc/profile.d/conda.csh"